# config/cliproxyapi.service
# CLIProxyAPI — systemd service (user)
#
# Installation:
#   mkdir -p ~/.config/systemd/user
#   cp cliproxyapi.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now cliproxyapi
#
# Check status:
#   systemctl --user status cliproxyapi
#
# View logs:
#   journalctl --user -u cliproxyapi -f

[Unit]
Description=CLIProxyAPI — OAuth to API proxy for RLM
After=network.target

[Service]
Type=simple
ExecStart=%h/.local/bin/cli-proxy-api --config %h/.cli-proxy-api/config.yaml
Restart=always
RestartSec=10
Environment=HOME=%h

[Install]
WantedBy=default.target
